# install/rdc-proxy.service
[Unit]
Description=RDC Local Proxy
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
# Kernel prerequisites
ExecStartPre=/sbin/modprobe br_netfilter
ExecStartPre=/bin/sh -c 'echo 1 > /proc/sys/net/bridge/bridge-nf-call-iptables'
# Route fwmarked traffic to local sockets (TPROXY plumbing)
ExecStartPre=/bin/sh -c 'ip rule add fwmark 1 lookup 100 2>/dev/null; ip route add local 0.0.0.0/0 dev lo table 100 2>/dev/null; true'
# TPROXY on the bridge: only divert RDC source -> cloud port
ExecStartPre=/bin/sh -c 'iptables -t mangle -C PREROUTING -i br0 -s ${RDC_IP} -p tcp --dport ${PROXY_PORT} -j TPROXY --tproxy-mark 0x1/0x1 --on-port ${PROXY_PORT} 2>/dev/null || iptables -t mangle -A PREROUTING -i br0 -s ${RDC_IP} -p tcp --dport ${PROXY_PORT} -j TPROXY --tproxy-mark 0x1/0x1 --on-port ${PROXY_PORT}'
# Broute: pull matching frames out of L2 bridging into L3 so TPROXY can intercept.
# ebtables default (nft shim) fails for broute redirect — ebtables-legacy works.
ExecStartPre=/bin/sh -c '/usr/sbin/ebtables-legacy -t broute -L BROUTING 2>/dev/null | grep -q "ip-src ${RDC_IP}.*ip-dport ${PROXY_PORT}" || /usr/sbin/ebtables-legacy -t broute -A BROUTING -p IPv4 --ip-proto tcp --ip-dport ${PROXY_PORT} --ip-src ${RDC_IP} -j redirect --redirect-target ACCEPT'

Environment="RDC_PROXY_CONFIG_DIR=/etc/rdc-proxy"
Environment="RDC_IP=10.0.0.50"
Environment="PROXY_PORT=5253"
EnvironmentFile=-/etc/rdc-proxy/rdc-proxy.env
ExecStart=/opt/rdc-proxy/venv/bin/python -m rdc_proxy

ExecStopPost=/bin/sh -c '/usr/sbin/ebtables-legacy -t broute -D BROUTING -p IPv4 --ip-proto tcp --ip-dport ${PROXY_PORT} --ip-src ${RDC_IP} -j redirect --redirect-target ACCEPT 2>/dev/null; true'
ExecStopPost=/bin/sh -c 'iptables -t mangle -D PREROUTING -i br0 -s ${RDC_IP} -p tcp --dport ${PROXY_PORT} -j TPROXY --tproxy-mark 0x1/0x1 --on-port ${PROXY_PORT} 2>/dev/null; true'
ExecStopPost=/bin/sh -c 'ip rule del fwmark 1 lookup 100 2>/dev/null; ip route del local 0.0.0.0/0 dev lo table 100 2>/dev/null; true'
ExecStopPost=/bin/sh -c 'echo 0 > /proc/sys/net/bridge/bridge-nf-call-iptables'

Restart=always
RestartSec=10s
User=root
Group=root
StandardOutput=journal
StandardError=journal
SyslogIdentifier=rdc-proxy

[Install]
WantedBy=multi-user.target
